# packaging.in/murphy-wait-for-launchpad-ready.path
[Path]
PathExists=/tmp/amd_ready
Unit=ico-homescreen.service
