# packaging.in/murphy-wait-for-launchpad-ready.path
[Path]
PathExists=/run/user/%U/amd_agent
Unit=ico-homescreen.service
