# packaging.in/murphy-wait-for-launchpad-ready.path
[Path]
PathExists=/run/user/%U/wayland-0
Unit=ico-homescreen.service
