# packaging.in/murphyd.service
[Unit]
Description=Murphy Resource Policy Daemon

[Service]
Type=simple
User=app
Group=app
ExecStart=/usr/bin/murphyd -t systemd -vvv -f
SmackExecLabel=User
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
